# etc/init.d/rtkit-openrc
#!/sbin/openrc-run

description="Realtime Policy and Watchdog Daemon"
command="/usr/libexec/rtkit-daemon"
supervisor=supervise-daemon

depend() {
    need dbus
    use logger
}
